# qr3_22.conf
[program:qr3_22]
command=python3 /home/curator3/qr3_22/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator3/qr3_22/test.err.log
stdout_logfile=/home/curator3/qr3_22/test.err.log
